# roles/gitlab/templates/gitlab-backup/cron.d/gitlab-backup.j2
MAILTO="{{ gitlab_backup_cron_email_to }}"

12 4 * * * root {{ gitlab_backup_program_install_dir }}/backup.sh
